# deploy/autobot.service
[Unit]
Description=Auto Bot Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=botuser
WorkingDirectory=/opt/auto_bot
EnvironmentFile=/etc/auto_bot.env
ExecStart=/opt/auto_bot/.venv/bin/python main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
